# templates/default/redis.upstart.conf.erb
description "Start the redis instance on port <%= @port %>"
author "Installed by chef redisio cookbook"

<% if @run_on_startup %>
start on runlevel [2345]
<% end %>
stop on runlevel [06]
limit nofile 200000 200000
setuid <%= @user %>
setgid <%= @group %>


<% if @shutdown_save or (@save and !@save.empty?) %>
pre-stop script
 echo "pre-stop script called"
 /usr/local/bin/redis-cli -p <%= @port %> save
 echo "finished saving"
end script
<% end %>



# If the job exits, restart it. Give up with more than 10 restarts in 30 seconds.
respawn
respawn limit 10 30

pre-start exec /usr/local/bin/redis-server <%= @configdir %>/<%= @port %>.conf && echo "starting redis-server"
post-stop exec /usr/local/bin/redis-cli -p <%= @port %> shutdown