# templates/default/redis@.service.erb
[Unit]
Description=Redis persistent key-value database
After=network.target

[Service]
ExecStart=<%= @bin_path %>/redis-server /etc/redis/%i.conf --daemonize no
ExecStop=/usr/bin/redis-shutdown
User=redis
Group=redis
LimitNOFILE=<%= @limit_nofile %>

[Install]
WantedBy=multi-user.target
