# templates/default/redis-server.conf.erb
description "redis server (<%= @name %>)"

start on runlevel [23]
stop on shutdown

<%- if @user && @group %>
setuid <%= @user %>
setgid <%= @group %>
<%- end %>

expect fork

<%= 'respawn' if @respawn %>

exec /usr/local/bin/redis-server /etc/redis/<%= @port %>.conf

